# pihole/systemd/restic-backup.timer
[Unit]
Description=Run restic backup every 24 hours

[Timer]
OnCalendar=daily
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
