# k3s-solo/update-motd.d/00-figlet
#!/bin/sh

. /etc/update-motd.d/quiet
. /etc/update-motd.d/colors

printf "\n"$GREEN
figlet "  "$(hostname -s)
printf $NONE
printf "\n"
